# sysutils/openrc/files/openrc.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: openrc
# REQUIRE: LOCAL
# BEFORE: LOGIN
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf[.local] to enable openrc
#
# openrc_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable openrc.

. /etc/rc.subr

name=openrc
rcvar=openrc_enable

start_cmd="${name}_start"
stop_cmd="${name}_stop"


openrc_start()
{
	/usr/local/sbin/openrc sysinit || exit 1
	/usr/local/sbin/openrc boot || exit 1
	/usr/local/sbin/openrc default

}

openrc_stop()
{
	LD_LIBRARY_PATH="/lib${LD_LIBRARY_PATH:+:}${LDLIBRARY_PATH}" ; export LD_LIBRARY_PATH
	action=${1:-shutdown}
	exec /usr/local/sbin/openrc "${action}"
}

load_rc_config $name

: ${openrc_enable="NO"}

run_rc_command "$1"
