# x11/kde4-workspace/pkg-plist
bin/startkde
bin/startkde4
bin/systemsettings
etc/dbus-1/system.d/org.kde.ksysguard.processlisthelper.conf
etc/rc.d/kdm4
include/kde4/kgreeterplugin.h
include/kde4/ksgrd/SensorAgent.h
include/kde4/ksgrd/SensorClient.h
include/kde4/ksgrd/SensorManager.h
include/kde4/ksgrd/SensorShellAgent.h
include/kde4/ksgrd/SensorSocketAgent.h
include/kde4/ksysguard/KTextEditVT.h
include/kde4/ksysguard/ProcessFilter.h
include/kde4/ksysguard/ProcessModel.h
include/kde4/ksysguard/ksignalplotter.h
include/kde4/ksysguard/ksysguardprocesslist.h
include/kde4/ksysguard/lsof.h
include/kde4/ksysguard/process.h
include/kde4/ksysguard/processes.h
include/kde4/kworkspace/kdisplaymanager.h
include/kde4/kworkspace/kwindowlistmenu.h
include/kde4/kworkspace/kworkspace.h
include/kde4/kworkspace/kworkspace_export.h
include/kde4/plasmaclock/calendar.h
include/kde4/plasmaclock/calendartable.h
include/kde4/plasmaclock/clockapplet.h
include/kde4/plasmaclock/plasmaclock_export.h
include/kde4/systemsettingsview/BaseData.h
include/kde4/systemsettingsview/BaseMode.h
include/kde4/systemsettingsview/MenuItem.h
include/kde4/systemsettingsview/MenuModel.h
include/kde4/systemsettingsview/MenuProxyModel.h
include/kde4/systemsettingsview/ModuleView.h
include/kde4/systemsettingsview/systemsettingsview_export.h
include/kde4/taskmanager/abstractgroupableitem.h
include/kde4/taskmanager/abstractgroupingstrategy.h
include/kde4/taskmanager/abstractsortingstrategy.h
include/kde4/taskmanager/groupmanager.h
include/kde4/taskmanager/launcheritem.h
include/kde4/taskmanager/startup.h
include/kde4/taskmanager/task.h
include/kde4/taskmanager/taskactions.h
include/kde4/taskmanager/taskgroup.h
include/kde4/taskmanager/taskitem.h
include/kde4/taskmanager/taskmanager.h
include/kde4/taskmanager/taskmanager_export.h
lib/cmake/KDE4Workspace/KDE4WorkspaceConfig.cmake
lib/cmake/KDE4Workspace/KDE4WorkspaceConfigVersion.cmake
lib/cmake/KDE4Workspace/KDE4WorkspaceLibraryTargets-%%CMAKE_BUILD_TYPE%%.cmake
lib/cmake/KDE4Workspace/KDE4WorkspaceLibraryTargets.cmake
lib/kde4/classic_mode.so
lib/kde4/icon_mode.so
lib/kde4/kded_kephal.so
lib/kde4/kgreet_classic.so
lib/kde4/kgreet_generic.so
lib/kde4/kgreet_winbind.so
lib/kde4/libexec/ksysguardprocesslist_helper
lib/kde4/plugins/designer/ksignalplotterwidgets.so
lib/kde4/plugins/designer/ksysguardlsofwidgets.so
lib/kde4/plugins/designer/ksysguardwidgets.so
lib/libkephal.so
lib/libkephal.so.4
lib/libkephal.so.%%KDE4_KDELIBS_VERSION%%
lib/libksgrd.so
lib/libksgrd.so.4
lib/libksgrd.so.%%KDE4_KDELIBS_VERSION%%
lib/libksignalplotter.so
lib/libksignalplotter.so.4
lib/libksignalplotter.so.%%KDE4_KDELIBS_VERSION%%
lib/libkworkspace.so
lib/libkworkspace.so.4
lib/libkworkspace.so.%%KDE4_KDELIBS_VERSION%%
lib/liblsofui.so
lib/liblsofui.so.4
lib/liblsofui.so.%%KDE4_KDELIBS_VERSION%%
lib/liboxygenstyle.so
lib/liboxygenstyle.so.4
lib/liboxygenstyle.so.%%KDE4_KDELIBS_VERSION%%
lib/liboxygenstyleconfig.so
lib/liboxygenstyleconfig.so.4
lib/liboxygenstyleconfig.so.%%KDE4_KDELIBS_VERSION%%
lib/libplasmaclock.so
lib/libplasmaclock.so.4
lib/libplasmaclock.so.%%KDE4_KDELIBS_VERSION%%
lib/libplasmagenericshell.so
lib/libplasmagenericshell.so.4
lib/libplasmagenericshell.so.%%KDE4_KDELIBS_VERSION%%
lib/libprocesscore.so
lib/libprocesscore.so.4
lib/libprocesscore.so.%%KDE4_KDELIBS_VERSION%%
lib/libprocessui.so
lib/libprocessui.so.4
lib/libprocessui.so.%%KDE4_KDELIBS_VERSION%%
lib/libsystemsettingsview.so
lib/libsystemsettingsview.so.2
lib/libtaskmanager.so
lib/libtaskmanager.so.4
lib/libtaskmanager.so.%%KDE4_KDELIBS_VERSION%%
share/applications/kde4/kdesystemsettings.desktop
share/applications/kde4/systemsettings.desktop
share/apps/cmake/modules/FindCkConnector.cmake
share/apps/cmake/modules/FindDBus.cmake
share/apps/cmake/modules/FindKephal.cmake
share/apps/cmake/modules/FindLibXKlavier.cmake
share/apps/cmake/modules/FindOpenGLES.cmake
share/apps/cmake/modules/FindPAM.cmake
share/apps/cmake/modules/FindSensors.cmake
share/apps/cmake/modules/PkgConfigGetVar.cmake
share/apps/cmake/modules/UnixAuth.cmake
share/apps/ksysguard/scripts/README
share/apps/ksysguard/scripts/smaps/helper.js
share/apps/ksysguard/scripts/smaps/index.html
share/apps/ksysguard/scripts/smaps/main.js
share/apps/ksysguard/scripts/smaps/smaps.desktop
share/apps/ksysguard/scripts/smaps/sorttable.js
share/apps/ksysguard/scripts/smaps/style.css
share/apps/plasma/packages/org.kde.desktop.widgetexplorer/contents/ui/AppletDelegate.qml
share/apps/plasma/packages/org.kde.desktop.widgetexplorer/contents/ui/Tooltip.qml
share/apps/plasma/packages/org.kde.desktop.widgetexplorer/contents/ui/main.qml
share/apps/plasma/packages/org.kde.desktop.widgetexplorer/metadata.desktop
share/apps/systemsettings/classic/main.html
share/apps/systemsettings/classic/systemsettings-classic.css
share/apps/systemsettings/systemsettings.kcfg
share/apps/systemsettings/systemsettingsui.rc
share/dbus-1/system-services/org.kde.ksysguard.processlisthelper.service
share/kde4/services/kded/kephal.desktop
share/kde4/services/settings-accessibility.desktop
share/kde4/services/settings-account-details.desktop
share/kde4/services/settings-application-and-system-notifications.desktop
share/kde4/services/settings-application-appearance-and-behavior.desktop
share/kde4/services/settings-application-appearance.desktop
share/kde4/services/settings-audio-and-video.desktop
share/kde4/services/settings-bluetooth.desktop
share/kde4/services/settings-classic-view.desktop
share/kde4/services/settings-desktop-appearance.desktop
share/kde4/services/settings-display.desktop
share/kde4/services/settings-hardware.desktop
share/kde4/services/settings-icon-view.desktop
share/kde4/services/settings-input-devices.desktop
share/kde4/services/settings-locale.desktop
share/kde4/services/settings-lost-and-found.desktop
share/kde4/services/settings-network-and-connectivity.desktop
share/kde4/services/settings-network-settings.desktop
share/kde4/services/settings-permissions.desktop
share/kde4/services/settings-personal-information.desktop
share/kde4/services/settings-power-management.desktop
share/kde4/services/settings-removable-devices.desktop
share/kde4/services/settings-sharing.desktop
share/kde4/services/settings-shortcuts-and-gestures.desktop
share/kde4/services/settings-startup-and-shutdown.desktop
share/kde4/services/settings-system-administration.desktop
share/kde4/services/settings-window-behaviour.desktop
share/kde4/services/settings-workspace-appearance-and-behavior.desktop
share/kde4/services/settings-workspace-behavior.desktop
share/kde4/servicetypes/plasma-layout-template.desktop
share/kde4/servicetypes/systemsettingscategory.desktop
share/kde4/servicetypes/systemsettingsexternalapp.desktop
share/kde4/servicetypes/systemsettingsview.desktop
share/polkit-1/actions/org.kde.ksysguard.processlisthelper.policy
