# sysutils/openrc/pkg-plist
bin/rc-status
etc/conf.d/bootmisc
etc/conf.d/fsck
etc/conf.d/hostname
etc/conf.d/ipfw
etc/conf.d/localmount
etc/conf.d/moused
etc/conf.d/netmount
etc/conf.d/network
etc/conf.d/powerd
etc/conf.d/rarpd
etc/conf.d/savecore
etc/conf.d/staticroute
etc/conf.d/syscons
etc/conf.d/tmpfiles
etc/conf.d/urandom
etc/init.d/adjkerntz
etc/init.d/bootmisc
etc/init.d/cron
etc/init.d/devd
etc/init.d/dumpon
etc/init.d/encswap
etc/init.d/fsck
etc/init.d/functions.sh
etc/init.d/hostid
etc/init.d/hostname
etc/init.d/ipfw
etc/init.d/local
etc/init.d/localmount
etc/init.d/loopback
etc/init.d/mixer
etc/init.d/moused
etc/init.d/netmount
etc/init.d/network
etc/init.d/newsyslog
etc/init.d/nscd
etc/init.d/osclock
etc/init.d/pf
etc/init.d/powerd
etc/init.d/rarpd
etc/init.d/rc-enabled
etc/init.d/root
etc/init.d/rpcbind
etc/init.d/s6-svscan
etc/init.d/savecache
etc/init.d/savecore
etc/init.d/staticroute
etc/init.d/swap
etc/init.d/swapfiles
etc/init.d/swclock
etc/init.d/syscons
etc/init.d/sysctl
etc/init.d/syslogd
etc/init.d/tmpfiles.setup
etc/init.d/urandom
etc/init.d/zfs
etc/init.d/zvol
etc/local.d/README
etc/rc
etc/rc.conf
etc/rc.d/openrc
etc/rc.devd
etc/rc.shutdown
etc/runlevels/boot/adjkerntz
etc/runlevels/boot/bootmisc
etc/runlevels/boot/dumpon
etc/runlevels/boot/fsck
etc/runlevels/boot/hostid
etc/runlevels/boot/hostname
etc/runlevels/boot/localmount
etc/runlevels/boot/loopback
etc/runlevels/boot/network
etc/runlevels/boot/newsyslog
etc/runlevels/boot/root
etc/runlevels/boot/savecore
etc/runlevels/boot/staticroute
etc/runlevels/boot/swap
etc/runlevels/boot/swapfiles
etc/runlevels/boot/syscons
etc/runlevels/boot/sysctl
etc/runlevels/boot/syslogd
etc/runlevels/boot/urandom
etc/runlevels/default/local
etc/runlevels/default/netmount
etc/runlevels/shutdown/savecache
etc/sysctl.d/README
include/einfo.h
include/rc.h
lib/libeinfo.so
lib/libeinfo.so.1
lib/librc.so
lib/librc.so.1
libexec/rc/bin/checkpath
libexec/rc/bin/ebegin
libexec/rc/bin/eend
libexec/rc/bin/eerror
libexec/rc/bin/eerrorn
libexec/rc/bin/eindent
libexec/rc/bin/einfo
libexec/rc/bin/einfon
libexec/rc/bin/eoutdent
libexec/rc/bin/esyslog
libexec/rc/bin/eval_ecolors
libexec/rc/bin/ewaitfile
libexec/rc/bin/ewarn
libexec/rc/bin/ewarnn
libexec/rc/bin/ewend
libexec/rc/bin/fstabinfo
libexec/rc/bin/get_options
libexec/rc/bin/is_newer_than
libexec/rc/bin/is_older_than
libexec/rc/bin/mountinfo
libexec/rc/bin/on_ac_power
libexec/rc/bin/rc-depend
libexec/rc/bin/save_options
libexec/rc/bin/service_crashed
libexec/rc/bin/service_get_value
libexec/rc/bin/service_hotplugged
libexec/rc/bin/service_inactive
libexec/rc/bin/service_set_value
libexec/rc/bin/service_started
libexec/rc/bin/service_started_daemon
libexec/rc/bin/service_starting
libexec/rc/bin/service_stopped
libexec/rc/bin/service_stopping
libexec/rc/bin/service_wasinactive
libexec/rc/bin/shell_var
libexec/rc/bin/vebegin
libexec/rc/bin/veend
libexec/rc/bin/veindent
libexec/rc/bin/veinfo
libexec/rc/bin/veoutdent
libexec/rc/bin/vewarn
libexec/rc/bin/vewend
libexec/rc/sbin/mark_service_failed
libexec/rc/sbin/mark_service_hotplugged
libexec/rc/sbin/mark_service_inactive
libexec/rc/sbin/mark_service_started
libexec/rc/sbin/mark_service_starting
libexec/rc/sbin/mark_service_stopped
libexec/rc/sbin/mark_service_stopping
libexec/rc/sbin/mark_service_wasinactive
libexec/rc/sbin/rc-abort
libexec/rc/sbin/swclock
libexec/rc/sh/functions.sh
libexec/rc/sh/gendepends.sh
libexec/rc/sh/init.sh
libexec/rc/sh/openrc-run.sh
libexec/rc/sh/rc-functions.sh
libexec/rc/sh/rc-mount.sh
libexec/rc/sh/s6.sh
libexec/rc/sh/start-stop-daemon.sh
libexec/rc/sh/tmpfiles.sh
libexec/rc/version
man/man3/ebegin.3.gz
man/man3/ebeginn.3.gz
man/man3/ebeginv.3.gz
man/man3/ebeginvn.3.gz
man/man3/ebracket.3.gz
man/man3/eend.3.gz
man/man3/eendv.3.gz
man/man3/eerror.3.gz
man/man3/eerrorn.3.gz
man/man3/eerrorx.3.gz
man/man3/eindent.3.gz
man/man3/eindentv.3.gz
man/man3/einfo.3.gz
man/man3/einfon.3.gz
man/man3/einfov.3.gz
man/man3/einfovn.3.gz
man/man3/eoutdent.3.gz
man/man3/eoutdentv.3.gz
man/man3/eprefix.3.gz
man/man3/ewarn.3.gz
man/man3/ewarnn.3.gz
man/man3/ewarnv.3.gz
man/man3/ewarnvn.3.gz
man/man3/ewarnx.3.gz
man/man3/ewend.3.gz
man/man3/ewendv.3.gz
man/man3/rc_config.3.gz
man/man3/rc_config_list.3.gz
man/man3/rc_config_load.3.gz
man/man3/rc_config_value.3.gz
man/man3/rc_deptree.3.gz
man/man3/rc_deptree_depend.3.gz
man/man3/rc_deptree_depends.3.gz
man/man3/rc_deptree_free.3.gz
man/man3/rc_deptree_load.3.gz
man/man3/rc_deptree_order.3.gz
man/man3/rc_deptree_update.3.gz
man/man3/rc_deptree_update_needed.3.gz
man/man3/rc_find_pids.3.gz
man/man3/rc_plugin_hook.3.gz
man/man3/rc_runlevel.3.gz
man/man3/rc_runlevel_exists.3.gz
man/man3/rc_runlevel_get.3.gz
man/man3/rc_runlevel_list.3.gz
man/man3/rc_runlevel_set.3.gz
man/man3/rc_runlevel_starting.3.gz
man/man3/rc_runlevel_stopping.3.gz
man/man3/rc_service.3.gz
man/man3/rc_service_add.3.gz
man/man3/rc_service_daemon_set.3.gz
man/man3/rc_service_daemons_crashed.3.gz
man/man3/rc_service_delete.3.gz
man/man3/rc_service_description.3.gz
man/man3/rc_service_exists.3.gz
man/man3/rc_service_extra_commands.3.gz
man/man3/rc_service_in_runlevel.3.gz
man/man3/rc_service_mark.3.gz
man/man3/rc_service_plugable.3.gz
man/man3/rc_service_resolve.3.gz
man/man3/rc_service_schedule_clear.3.gz
man/man3/rc_service_schedule_start.3.gz
man/man3/rc_service_started_daemon.3.gz
man/man3/rc_service_state.3.gz
man/man3/rc_service_value_get.3.gz
man/man3/rc_service_value_set.3.gz
man/man3/rc_services_in_runlevel.3.gz
man/man3/rc_services_in_state.3.gz
man/man3/rc_services_scheduled.3.gz
man/man3/rc_services_scheduled_by.3.gz
man/man3/rc_stringlist.3.gz
man/man3/rc_stringlist_add.3.gz
man/man3/rc_stringlist_addu.3.gz
man/man3/rc_stringlist_delete.3.gz
man/man3/rc_stringlist_free.3.gz
man/man3/rc_stringlist_new.3.gz
man/man3/rc_stringlist_sort.3.gz
man/man3/rc_yesno.3.gz
man/man8/openrc-run.8.gz
man/man8/openrc.8.gz
man/man8/rc-service.8.gz
man/man8/rc-status.8.gz
man/man8/rc-update.8.gz
man/man8/service.8.gz
man/man8/start-stop-daemon.8.gz
sbin/init
sbin/openrc
sbin/openrc-run
sbin/rc
sbin/rc-service
sbin/rc-update
sbin/runscript
sbin/service
sbin/start-stop-daemon
@dir etc/runlevels/sysinit
@dir libexec/rc/init.d
@dir libexec/rc/tmp
@dir /usr/lib/pkgconfig
